# scripts/jf-sebastian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd plist for J.F. Sebastian on macOS.

  Install:
    1. Replace ALL __EDIT_ME__ placeholders with absolute paths and your username.
       (sed one-liner: sed -i '' "s|__EDIT_ME_REPO__|$PWD|g; s|__EDIT_ME_USER__|$USER|g" scripts/jf-sebastian.plist)
    2. Copy to ~/Library/LaunchAgents/com.jf-sebastian.supervisor.plist
    3. launchctl bootstrap gui/$UID ~/Library/LaunchAgents/com.jf-sebastian.supervisor.plist

  KeepAlive is set so launchd restarts the *supervisor itself* if it crashes.
  The supervisor in turn restarts the actual JF Sebastian child on its crashes,
  so this is defense in depth: launchd watches the supervisor; supervisor
  watches the child.

  The supervisor manages its own log rotation (see scripts/supervisor.py).
  We do NOT set StandardOutPath/StandardErrorPath here — the supervisor's
  RotatingFileHandler writes to SUPERVISOR_LOG_PATH directly.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jf-sebastian.supervisor</string>

    <key>ProgramArguments</key>
    <array>
        <string>__EDIT_ME_REPO__/venv/bin/python</string>
        <string>__EDIT_ME_REPO__/scripts/supervisor.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__EDIT_ME_REPO__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HEARTBEAT_FILE</key>
        <string>/tmp/jf_sebastian.heartbeat</string>
        <key>CRASH_REPORT_DIR</key>
        <string>__EDIT_ME_REPO__/crash_reports</string>
        <key>LOG_PATH</key>
        <string>__EDIT_ME_REPO__/jf_sebastian.log</string>
        <key>SUPERVISOR_LOG_PATH</key>
        <string>__EDIT_ME_REPO__/supervisor.log</string>
        <!-- ffmpeg lives in /opt/homebrew/bin (Apple Silicon) or /usr/local/bin (Intel) -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/__EDIT_ME_USER__</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>LC_ALL</key>
        <string>en_US.UTF-8</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <!-- Reduce CUDA allocator fragmentation on Jetson; ignored on Mac (MPS/CPU). -->
        <key>PYTORCH_CUDA_ALLOC_CONF</key>
        <string>expandable_segments:True</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <!-- Restart the supervisor itself if it ever exits abnormally. -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Don't relaunch a wedged supervisor more than once every 30s. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Audio device access requires the Aqua session. -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
