# scripts/jf-sebastian.service
# systemd user unit for J.F. Sebastian.
#
# Install:
#   1. Edit %h placeholders to point at your install (or keep them for ~/jf-sebastian).
#      Note: ExecStart assumes a venv at <install>/venv. For a system Python install
#      or a different venv location, edit the python path below.
#   2. Copy to ~/.config/systemd/user/jf-sebastian.service
#   3. systemctl --user daemon-reload
#      systemctl --user enable --now jf-sebastian.service
#
# The supervisor handles restart/backoff/watchdog of the JF Sebastian child.
# systemd's Restart=on-failure here is just a safety net for the supervisor itself
# crashing — defense in depth.

[Unit]
Description=J.F. Sebastian animatronic conversation system (supervised)
After=network-online.target pipewire-pulse.service
Wants=pipewire-pulse.service
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
WorkingDirectory=%h/jf-sebastian
ExecStart=%h/jf-sebastian/venv/bin/python %h/jf-sebastian/scripts/supervisor.py
Environment=HEARTBEAT_FILE=%t/jf_sebastian.heartbeat
Environment=CRASH_REPORT_DIR=%h/jf-sebastian/crash_reports
Environment=LOG_PATH=%h/jf-sebastian/jf_sebastian.log
Environment=SUPERVISOR_LOG_PATH=%h/jf-sebastian/supervisor.log
Environment=PYTHONUNBUFFERED=1
# Reduce CUDA caching-allocator fragmentation on Jetson during long RVC sessions.
# Ignored on non-CUDA platforms.
Environment=PYTORCH_CUDA_ALLOC_CONF=expandable_segments:True
# Locale: systemd user units inherit very little; explicit is safer for ffmpeg
# subprocess output and any non-ASCII personality strings.
Environment=LANG=en_US.UTF-8
Environment=LC_ALL=en_US.UTF-8
Restart=on-failure
RestartSec=10
# Generous TimeoutStopSec — supervisor's own SHUTDOWN_GRACE_SECS (default 10s)
# plus interpreter teardown plus in-flight TTS/RVC cleanup needs more than 15s.
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jf-sebastian

[Install]
WantedBy=default.target
